# bin/crontab
* * * * * ssh-agent $HOME/bin/duc.sh >/dev/null 2>&1
